# deploy/oxiterm.service
[Unit]
Description=OxiTerm Weather Demo Server
After=network.target

[Service]
Type=simple
User=slavekm
WorkingDirectory=/my_data/oxiterm
Environment=OXITERM_PASSWORD=krakow
Environment=RUST_LOG=info
ExecStart=/my_data/oxiterm/target/release/oxiterm-server
Restart=always
RestartSec=5
StandardOutput=append:/tmp/oxiterm.log
StandardError=append:/tmp/oxiterm.log

[Install]
WantedBy=multi-user.target
